# deploy/launchd/com.ourworlds.quant.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ourworlds.quant.app</string>

  <key>WorkingDirectory</key>
  <string>/tmp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Volumes/EXTDISK/QUANT/QUANT/ourworld-quant &amp;&amp; exec deploy/run-public-app.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OWQ_ENV_FILE</key>
    <string>/Volumes/EXTDISK/QUANT/QUANT/ourworld-quant/deploy/public.env</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Volumes/EXTDISK/QUANT/QUANT/ourworld-quant/data/logs/app.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Volumes/EXTDISK/QUANT/QUANT/ourworld-quant/data/logs/app.err.log</string>
</dict>
</plist>
